# components/master/launchd/com.subctl.master.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  subctl master (subctl master daemon) launchd template.

  Template — placeholders are substituted at install time by
  lib/master.sh:subctl_master_enable():
    __OWNER__      reverse-DNS prefix (e.g. com.subctl)
    __BUN__        absolute path to the bun binary
    __SERVER_TS__  absolute path to components/master/server.ts
    __HOME__       absolute path to user's $HOME

  Boot expectations from the daemon (server.ts):
    - emit a heartbeat log line within 30s of boot. server.ts logs
      `[master] booting subctl master ...` immediately and a heartbeat every
      60s thereafter; ThrottleInterval=10 keeps successive crash-restarts
      from runaway-detecting.

  Both stdout and stderr stream to the SAME file
  (~/Library/Logs/subctl/master.log) — server.ts uses console.error
  for all structured logs (so a single file gives you the full picture
  without merging streams).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__OWNER__.master</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BUN__</string>
        <string>run</string>
        <string>__SERVER_TS__</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Conditional KeepAlive — restart on crash (Crashed:true) or non-zero
         exit (SuccessfulExit:false), but DO NOT restart after a clean
         operator-initiated stop. Without this, `subctl master disable`
         (which sends SIGTERM and exits 0) would still get respawned.
         Matches the dashboard plist's pattern. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- ThrottleInterval 30s (was 10). Today's death spiral 2026-05-10: LM
         Studio crashed, master ctx-pin hung 60s, daemon exited, launchd
         relaunched within 10s, master hung 60s again, repeat. Macos's
         internal respawn-limit detector flagged the job as failing and
         refused further restarts. 30s gives the environment more breathing
         room between crash-attempts and resets launchd's failure counter
         more aggressively. Combined with v2.5.3's 2s LM Studio reachability
         probe, the daemon now boots fast even when LM Studio is dead. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- ExitTimeOut 20s: launchd sends SIGTERM, waits up to 20s for graceful
         shutdown (transcript flush, SSE drain), then SIGKILL. Prevents
         zombie processes if a daemon's shutdown path hangs. -->
    <key>ExitTimeOut</key>
    <integer>20</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/subctl/master.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/subctl/master.log</string>
</dict>
</plist>
